# meta-agl-demo/recipes-demo/agl-persistent-storage-api/files/agl-persistent-storage-api.service
[Unit]
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/agl-service-persistent-storage
Restart=on-failure
User=agl-driver

[Install]
WantedBy=default.target
