# meta-agl-demo/recipes-demo/agl-service-hvac/agl-service-hvac/agl-service-hvac.service
[Service]
Type=simple
ExecStart=/usr/sbin/agl-service-hvac
Restart=on-failure

[Install]
WantedBy=default.target
