# meta-agl-demo/recipes-demo/agl-service-audiomixer/agl-service-audiomixer/agl-service-audiomixer.service
[Unit]
Requires=pipewire.service
After=pipewire.service

[Service]
Type=simple
ExecStart=/usr/sbin/agl-service-audiomixer
Restart=on-failure

[Install]
WantedBy=default.target
